# deploy/wetokenbot.ini
[program:wetokenbot]
command = /usr/bin/python run.py
directory=/opt/WeTokenBot/wetokebot/
user = root
redirect_stderr = root
stdout_logfile = /opt/WeTokenBot/run.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopasgroup=true
startretries = 3